# j-entrypoint/cron-shell/crontab
* 4 * * * /j-entrypoint/cron-shell/backup.sh
